# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[program:web-api]
command=/usr/local/bin/python3 /opt/tracker-autovisit/web-api.py
directory=/opt/tracker-autovisit
autostart=true
autorestart=true
startretries=10
environment=HOME="/root",PATH="/usr/local/bin:/usr/bin:/bin"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=30
