# systemd/autovisit-web.service
[Unit]
Description=tracker-autovisit web API (gestion des sites + parametres + auth)
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/tracker-autovisit/web-api.py
Environment=MAL_CT_IP=${CT_IP}
Environment=MAL_LAN_CIDR=${LAN_CIDR}
Environment=MAL_WEBROOT=/var/www/autovisit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
